# etc/s6-overlay/s6-rc.d/app/run
#!/bin/sh
echo "[s6-app] === RUN SCRIPT EXECUTED ==="
exec cross-env NODE_MODE=all node index.js
