# etc/s6-overlay/s6-rc.d/app/run
#!/bin/sh
echo "[s6-app] Starting Node.js application..."
exec node index.js
