# etc/s6-overlay/s6-rc.d/app/run
#!/bin/sh
echo "[s6-app] === RUN SCRIPT EXECUTED ==="
# 使用原生环境变量赋值，避免依赖 devDependencies 中的 cross-env
# 堆内存上限可通过 MAX_HEAP 环境变量配置，默认 512MB
MAX_HEAP="${MAX_HEAP:-512}"
NODE_OPTIONS="${NODE_OPTIONS:-} --max-old-space-size=${MAX_HEAP} --max-semi-space-size=16 --expose-gc"

# 条件加载 OpenTelemetry SDK（仅在配置了 New Relic License Key 时）
# --import 必须在业务模块加载前执行，以拦截 http/ioredis/undici 模块
# 默认采样率 10%，避免 256MB 容器在高负载下 span 队列溢出及 NR ingest 成本过高
if [ -n "${NEW_RELIC_LICENSE_KEY:-}" ]; then
  NODE_OPTIONS="$NODE_OPTIONS --import=${APP_DIR:-/app}/src/telemetry/tracing.js"
  OTEL_TRACES_SAMPLER="${OTEL_TRACES_SAMPLER:-parentbased_traceidratio}"
  OTEL_TRACES_SAMPLER_ARG="${OTEL_TRACES_SAMPLER_ARG:-0.1}"
  export OTEL_TRACES_SAMPLER OTEL_TRACES_SAMPLER_ARG
fi

NODE_MODE=all exec node index.js
