# etc/s6-overlay/s6-rc.d/app/run
#!/bin/sh
echo "[s6-app] === RUN SCRIPT EXECUTED ==="
# 使用原生环境变量赋值，避免依赖 devDependencies 中的 cross-env
# 设置 Node.js 堆内存上限：容器仅 256MB，Node 默认堆约 256MB 会吃光内存
# 限制到 200MB 给 OS/rclone/cloudflared 留空间，并启用高效 GC
NODE_OPTIONS="${NODE_OPTIONS:-} --max-old-space-size=200 --max-semi-space-size=16 --expose-gc"
NODE_MODE=all exec node index.js
