# etc/s6-overlay/s6-rc.d/app/run
#!/bin/sh
echo "[s6-app] === RUN SCRIPT EXECUTED ==="
# 使用原生环境变量赋值，避免依赖 devDependencies 中的 cross-env
NODE_MODE=all exec node index.js
