# etc/s6-overlay/s6-rc.d/cloudflared/run
#!/bin/sh
# 检查是否启用tunnel
if [ "$TUNNEL_ENABLED" != "true" ]; then
    echo "[s6-cloudflared] Tunnel disabled (TUNNEL_ENABLED=$TUNNEL_ENABLED)"
    echo "[s6-cloudflared] Exiting successfully (not needed)"
    # 立即退出，让s6知道这个服务"完成"了（不是失败）
    exit 0
fi

PORT=${PORT:-7860}
TUNNEL_METRICS_PORT=${TUNNEL_METRICS_PORT:-2000}
echo "[s6-cloudflared] Starting tunnel to 127.0.0.1:$PORT with metrics on 127.0.0.1:$TUNNEL_METRICS_PORT"

# 启动cloudflared并提取临时URL写入文件
cloudflared tunnel --url http://127.0.0.1:$PORT \
    --metrics 127.0.0.1:$TUNNEL_METRICS_PORT \
    --no-autoupdate > /tmp/cloudflared.log 2>&1 &

# 异步提取 URL
(
    timeout=60
    while [ $timeout -gt 0 ]; do
        if [ -f /tmp/cloudflared.log ]; then
            URL=$(grep -o 'https://[-0-9a-z]*\.trycloudflare\.com' /tmp/cloudflared.log | head -n 1)
            if [ -n "$URL" ]; then
                echo "$URL" > /tmp/cloudflared.url
                echo "[s6-cloudflared] Captured quick tunnel URL: $URL"
                break
            fi
        fi
        sleep 2
        timeout=$((timeout - 2))
    done
) &
wait

