# etc/s6-overlay/s6-rc.d/cloudflared/run
#!/bin/sh
# 检查是否启用tunnel
if [ "$TUNNEL_ENABLED" != "true" ]; then
    echo "[s6-cloudflared] Tunnel disabled (TUNNEL_ENABLED=$TUNNEL_ENABLED)"
    echo "[s6-cloudflared] Exiting successfully (not needed)"
    # 立即退出，让s6知道这个服务"完成"了（不是失败）
    exit 0
fi

PORT=${PORT:-7860}
TUNNEL_METRICS_PORT=${TUNNEL_METRICS_PORT:-2000}
echo "[s6-cloudflared] Starting tunnel to localhost:$PORT with metrics on $TUNNEL_METRICS_PORT"

# 启动cloudflared并提取临时URL写入文件
# 使用 tee 将日志同时输出到标准错误（供s6日志查看）和管道
cloudflared tunnel --url http://localhost:$PORT \
    --metrics localhost:$TUNNEL_METRICS_PORT \
    --no-autoupdate 2>&1 | tee /dev/stderr | while read -r line; do
    # 匹配 trycloudflare.com 域名
    URL=$(echo "$line" | grep -o 'https://[-0-9a-z]*\.trycloudflare\.com')
    if [ -n "$URL" ]; then
        echo "$URL" > /tmp/cloudflared.url
        echo "[s6-cloudflared] Captured quick tunnel URL: $URL"
    fi
done
