# etc/s6-overlay/s6-rc.d/cloudflared/run
#!/bin/sh
echo "[s6-cloudflared] === RUN SCRIPT EXECUTED ===" # Debug: Confirm script execution
# 由 Node.js 通过 s6-svc -u 触发启动
# 这样可以确保 Infisical 加载的变量生效

PORT=${PORT:-7860}
TUNNEL_METRICS_PORT=${TUNNEL_METRICS_PORT:-2000}
echo "[s6-cloudflared] Starting tunnel to 127.0.0.1:$PORT with metrics on 127.0.0.1:$TUNNEL_METRICS_PORT"

# 确保 /tmp 目录存在
mkdir -p /tmp

# 启动 cloudflared
# 注意：s6 期望 run 脚本在前台运行以监控进程
echo "[s6-cloudflared] Starting cloudflared tunnel..."

# 启动异步提取 URL 任务
(
    echo "[s6-cloudflared] URL extractor started, waiting for log file..."
    timeout=60
    while [ $timeout -gt 0 ]; do
        if [ -f /tmp/cloudflared.log ]; then
            URL=$(grep -oE 'https://[a-zA-Z0-9.-]+\.trycloudflare\.com' /tmp/cloudflared.log | head -n 1)
            if [ -n "$URL" ]; then
                echo "$URL" > /tmp/cloudflared.url
                echo "[s6-cloudflared] SUCCESS: Captured quick tunnel URL: $URL"
                break
            fi
        fi
        sleep 2
        timeout=$((timeout - 2))
    done
    if [ $timeout -le 0 ]; then
        echo "[s6-cloudflared] ERROR: Timed out waiting for tunnel URL in /tmp/cloudflared.log"
    fi
) &

# 前台启动 cloudflared，使用 tee 将输出同时写入文件和 stdout/stderr
cloudflared tunnel --url "http://127.0.0.1:$PORT" \
    --metrics "127.0.0.1:$TUNNEL_METRICS_PORT" \
    --no-autoupdate 2>&1 | tee /tmp/cloudflared.log


